# ecui-web-small-teststand.service
Description=Service for ECUI Webserver

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=/home/pi/.nvm/versions/node/v12.10.0/bin/node /home/pi/web_ecui_houbolt/server.js --smallTeststand
WorkingDirectory=/home/pi/web_ecui_houbolt/
Restart=always
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target